# scripts/claude-metrics-proxy/com.cyborg.claude-metrics-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cyborg.claude-metrics-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/Caskroom/miniconda/base/bin/python</string>
        <string>proxy.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/victor.peng/code/cyborg/scripts/claude-metrics-proxy</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-metrics-proxy.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-metrics-proxy.err</string>
</dict>
</plist>
